# taskdaemon.ini
[program: Task]
directory=/home/auser/TRClient/Agent/OpenLab
command=/home/auser/TRClient/Agent/OpenLab/venv/bin/openlab_agent task -f configs/agent.yaml --vhost Development -l debug
autostart=false
stopasgroup=true
killasgroup=true
autorestart=true
startsecs=5
startretries=3
stdout_logfile=/home/auser/TRClient/Agent/OpenLab/logs/taskrecv.log
redirect_stderr = true
stdout_logfile_maxbytes = 300MB
stdout_logfile_backups = 2
user=auser
